# ark2020.conf
[program:ark2020]
command = gunicorn backend.wsgi_prod -b 127.0.0.1:8000
directory = /root/ark2020/ark/backend
autostart = true
autorestart = true