# ark2020.conf
[program:ark2020]
# Use full path to gunicorn from virtual environment
# If venv is in /root/ark2020/ark/venv/, update the path accordingly
command = /root/ark2020/ark/venv/bin/gunicorn backend.wsgi_prod -b 127.0.0.1:8888 -w 12
directory = /root/ark2020/ark/backend/
autostart = true
autorestart = true
# Set environment variables
environment=PATH="/root/ark2020/ark/venv/bin:%(ENV_PATH)s"
# User to run as (adjust if needed)
user = root
# Logging
stdout_logfile = /var/log/ark2020/access.log
stderr_logfile = /var/log/ark2020/error.log